# vision_sdk/linux/targetfs/var/lib/opkg/info/busybox.list
/bin/sh
/bin/busybox.suid
/etc/init.d/inetd.busybox
/bin/busybox
/etc/busybox.links.nosuid
/etc/busybox.links.suid
/etc/inetd.conf
/bin/busybox.nosuid
