# vision_sdk/linux/targetfs/var/lib/opkg/info/nfs-utils-client.list
/var/lib/nfs/rmtab
/sbin/mount.nfs4
/usr/sbin/showmount
/etc/nfs-utils.conf
/var/lib/nfs/sm/
/usr/sbin/rpc.idmapd
/sbin/umount.nfs4
/etc/init.d/nfscommon
/var/lib/nfs/state
/sbin/mount.nfs
/var/lib/nfs/v4recovery/
/usr/sbin/start-statd
/var/lib/nfs/xtab
/var/lib/nfs/
/var/lib/nfs/statd/
/usr/sbin/nfsstat
/var/lib/nfs/etab
/sbin/umount.nfs
/var/lib/nfs/sm.bak/
/usr/sbin/rpc.statd
/usr/sbin/sm-notify
