# vision_sdk/linux/targetfs/var/lib/opkg/info/nfs-utils.list
/etc/init.d/nfsserver
/usr/sbin/exportfs
/usr/sbin/rpc.nfsd
/usr/sbin/rpc.mountd
/sbin/osd_login
