# vision_sdk/linux/targetfs/var/lib/opkg/info/rng-tools.list
/etc/init.d/rng-tools
/usr/bin/rngtest
/etc/default/default
/usr/sbin/rngd
