# vision_sdk/linux/targetfs/var/lib/opkg/info/rpcbind.list
/usr/sbin/rpcbind
/etc/rpcbind.conf
/usr/sbin/rpcinfo
/etc/init.d/rpcbind
