# vision_sdk/linux/targetfs/var/lib/opkg/info/run-postinsts.list
/etc/init.d/run-postinsts
/usr/sbin/run-postinsts
