# vision_sdk/linux/targetfs/var/lib/opkg/info/udev-cache.list
/etc/init.d/udev-cache
/etc/default/udev-cache
